# Linux/Linux学习笔记/shell/Chapter6_文件归档及加密等/etc/systemd/system/snap-core-3748.mount
[Unit]
Description=Mount unit for core
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/core_3748.snap
Where=/snap/core/3748
Type=squashfs
Options=nodev,ro

[Install]
WantedBy=multi-user.target
